# bmc/usr/lib/systemd/system/sonic-usb-network-init.service.d/10-hw-management-bmc.conf
# SPDX-FileCopyrightText: NVIDIA CORPORATION & AFFILIATES
#
# SONiC's sonic-usb-network-init runs dhclient on usb0. hw-management-bmc configures
# the same interface with a static Address= via systemd-networkd (see
# /etc/systemd/network/00-hw-management-bmc-usb0.network). Running both breaks
# addressing; AppArmor may also block dhclient on this interface.
#
# When our .network file is present (rendered at boot by hw-management-bmc-plat-specific-preps),
# skip the SONiC one-shot. If usb0 .network generation is skipped (invalid USB0_ADDRESS),
# this condition passes and SONiC can still run dhclient.
[Unit]
ConditionPathExists=!/etc/systemd/network/00-hw-management-bmc-usb0.network
