# bmc/usr/lib/systemd/system/sonic-usb-network-init.service.d/10-hw-management-bmc.conf
# SPDX-FileCopyrightText: NVIDIA CORPORATION & AFFILIATES
#
# SONiC's sonic-usb-network-init runs dhclient on usb0. hw-management-bmc configures
# the same interface with a static Address= via systemd-networkd (see
# /etc/systemd/network/00-hw-management-bmc-usb0.network). Running both breaks
# addressing; AppArmor may also block dhclient on this interface.
#
# When our .network file is present (rendered at boot by hw-management-bmc-plat-specific-preps),
# this unit is skipped (ConditionPathExists=! fails). When USB0_MANAGED_BY_NOS=1,
# plat-specific-preps does not write .network, so the condition passes and
# sonic-usb-network-init may run. Same if .network generation is skipped
# (invalid USB0_ADDRESS).
[Unit]
ConditionPathExists=!/etc/systemd/network/00-hw-management-bmc-usb0.network
