# debian/hw-management.hw-management-tc.service
[Unit]
Description=Thermal control service (ver 2.0) of Nvidia systems
After=hw-management.service
Requires=hw-management.service
PartOf=hw-management.service
Documentation=man:hw-management-tc.service(8)

StartLimitIntervalSec=1200
StartLimitBurst=15

[Service]
ExecStart=/bin/sh -c "/usr/bin/hw_management_thermal_control.py"
ExecStop=/bin/kill $MAINPID
TimeoutStopSec=1

Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
