# debian/postinst
#!/bin/sh
set -e

if [ ! -L /sbin/init ]; then
	#initd detected
	cp -f /usr/bin/hw-management.sh /etc/init.d/hw-management
	if [ -x "/etc/init.d/hw-management" ]; then
		update-rc.d hw-management defaults >/dev/null
		invoke-rc.d hw-management start || exit $?
	fi
fi
