# bmc/usr/lib/systemd/system/hw-management-bmc-early-i2c-init.service
# SPDX-FileCopyrightText: NVIDIA CORPORATION & AFFILIATES
[Unit]
Description=BMC Early I2C Device Initialization
DefaultDependencies=no
# Requires= (in addition to After=) pulls early-config into the same
# transaction on ad-hoc starts so the JSON config it deploys is on disk
# before this unit reads it.  See issue #4992267.
Requires=hw-management-bmc-early-config.service
After=hw-management-bmc-early-config.service
Before=hw-management-bmc-sync-ethaddr.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/hw-management-bmc-early-i2c-init.sh
SyslogIdentifier=hw-management-bmc-early-i2c-init

[Install]
WantedBy=sysinit.target
